# systemd_service/transcription_app.service
# systemd_service/transcription_app.service
[Unit]
Description=Transcription Flask App
After=network.target sound.target

[Service]
User=root
WorkingDirectory=/opt/transcript/
ExecStart=/usr/local/bin/gunicorn --workers 3 --bind 0.0.0.0:5000 app:app

Restart=always

[Install]
WantedBy=multi-user.target
