# systemd_service/transcription_app.service
# systemd_service/transcription_app.service
[Unit]
Description=Transcription Flask App
After=network.target

[Service]
User=transcript
WorkingDirectory=/home/transcript/git/transcript/
Environment="PATH=/home/transcript/git/transcript/transcription_env/bin"
ExecStart=/home/transcript/git/transcript/transcription_env/bin/gunicorn --workers 3 --bind 0.0.0.0:5000 app:app

Restart=always

[Install]
WantedBy=multi-user.target
